# docker/supervisor.conf
[program:serviceweb]
command = /usr/local/bin/uwsgi --ini /app/uwsgi-serviceweb.ini

[program:nginx]
command = /usr/sbin/nginx
